# marlow-compositor.service
[Unit]
Description=Marlow Compositor
# No After/Before graphical-session — WE are the graphical session

[Service]
Type=simple
ExecStart=/home/josemarlow/marlow-compositor/target/release/marlow-compositor
Environment=XDG_RUNTIME_DIR=/run/user/1000
Environment=XDG_SESSION_TYPE=wayland
Environment=WAYLAND_DISPLAY=wayland-1
Environment=PATH=/home/josemarlow/.local/bin:/home/josemarlow/.cargo/bin:/usr/local/bin:/usr/bin:/bin
Restart=on-failure
RestartSec=3

[Install]
WantedBy=default.target
