# session/marlow-mako.service
[Unit]
Description=Marlow notification daemon (mako)
PartOf=marlow-session.target

[Service]
Type=simple
ExecStart=/usr/bin/mako
Restart=on-failure
RestartSec=3

[Install]
WantedBy=marlow-session.target
