# session/marlow-waybar.service
[Unit]
Description=Marlow top panel (waybar)
PartOf=marlow-session.target

[Service]
Type=simple
ExecStart=/usr/bin/waybar -c /home/josemarlow/marlow-compositor/config/waybar/config.jsonc -s /home/josemarlow/marlow-compositor/config/waybar/style.css
Restart=on-failure
RestartSec=3

[Install]
WantedBy=marlow-session.target
