# session/marlow-daemon.service
[Unit]
Description=Marlow AI agent daemon
After=marlow-session.target
PartOf=marlow-session.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 -c "from marlow.daemon_linux import main; main()"
WorkingDirectory=/home/josemarlow/marlow
Environment=PATH=/home/josemarlow/.local/bin:/usr/local/bin:/usr/bin:/bin
EnvironmentFile=-/home/josemarlow/.config/marlow/env
Restart=on-failure
RestartSec=3

[Install]
WantedBy=marlow-session.target
